# src/etc/supervisor/conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -display :0 -forever -nopw -quiet -listen 0.0.0.0 -xkb
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
